# arch/streamplace.service
[Unit]
Description=Live video for the AT Protocol. Solving video for everybody forever.
Documentation=https://stream.place/docs
After=network-online.target
Wants=network-online.target
Requires=network.target

[Service]
Type=exec
ExecStart=/usr/bin/streamplace
Restart=always
RestartSec=10
TimeoutStartSec=60
TimeoutStopSec=30

# Environment
Environment=NODE_ENV=production
Environment=STREAMPLACE_CONFIG_DIR=/etc/streamplace
Environment=STREAMPLACE_DATA_DIR=/var/lib/streamplace
Environment=STREAMPLACE_LOG_LEVEL=info
Environment=STREAMPLACE_HTTP_PORT=8080
Environment=STREAMPLACE_HTTPS_PORT=8443

# Additional environment file for user customization
EnvironmentFile=-/etc/streamplace/environment

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RemoveIPC=yes
PrivateTmp=yes

# Allow access to necessary directories
ReadWritePaths=/var/lib/streamplace
ReadOnlyPaths=/etc/streamplace

# Network access
PrivateNetwork=no
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=2G
CPUQuota=200%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=streamplace

# Working directory
WorkingDirectory=/var/lib/streamplace

[Install]
WantedBy=multi-user.target
